# packer/alfresco/cron.d/alfresco
SHELL=/bin/bash
PATH=/bin:/usr/bin

# m h  dom mon dow   command
@reboot     alfresco    /home/ubuntu/setup-alfresco.sh
